# supervisord.conf
[program:api]
command = node milibet.js
directory = /home/dotcloud/current
